# shipyard/docs/ops/templates/linux-mission/shipyard-vercel-sync.timer
[Unit]
Description=Run Shipyard backup Vercel deploy sync every two hours

[Timer]
OnBootSec=10m
OnUnitActiveSec=2h
Persistent=true
Unit=shipyard-vercel-sync.service

[Install]
WantedBy=timers.target
